# examples/rc.local.example
#!/bin/bash
ip addr add fe80::1/64 dev eth0 scope link
ip addr add 2001:db8:50::1/64 dev eth0
ip -6 route add default via 2001:db8:50::180 dev eth0
systemctl restart kea-dhcp6
